# system/etc/udev/rules.d/10-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/devon/.local/bin/runx /home/devon/.local/bin/msu -s 1"
