# system/etc/systemd/system/compute-batt-level.service
[Unit]
Description=Computes a combined battery level

[Service]
Type=oneshot
ExecStart=/home/devon/.local/bin/compute-combined-battery-level
